# src/filesystem/configs/etc/rustysd/services/default.target
[Unit]
Description= Spawn all the needed system services