# src/filesystem/configs/etc/rustysd/services/syslogd.service
[Unit]
Description= Start syslogd

[Service]
ExecStart=/sbin/syslogd -n
Type=simple
Restart=always

[Install]
WantedBy=default.target