# src/filesystem/configs/etc/rustysd/services/udhcpcd.service
[Unit]
Description= Start udhcpcd to get an ip

[Service]
ExecStart=/sbin/udhcpc -f -S
Type=simple
Restart=always

[Install]
WantedBy=default.target